# infra/systemd/chesterwc-backup.service
[Unit]
Description=Daily SQLite backup for chesterwindowcleaner

[Service]
Type=oneshot
ExecStart=/opt/chesterwc/infra/scripts/backup.sh
User=chesterwc
Group=chesterwc
